# misc/deb/DEBIAN/postinst
#!/bin/bash

cd /usr/share/ezpaarse/
make

update-rc.d ezpaarse defaults

/etc/init.d/ezpaarse start